# server/cron/field-prewarm.timer
# /etc/systemd/system/field-prewarm.timer
[Unit]
Description=Run matrix prewarm nightly

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target